# petalinux_project/project-spec/meta-user/recipes-apps/fractal-explorer/files/init
#!/bin/sh
### BEGIN INIT INFO
# Provides:       fractal-explorer
# Required-Start: $local_fs $remote_fs
# Required-Stop:  $local_fs $remote_fs
# Default-Start:  2 3 4 5
# Default-Stop:   0 1 6
### END INIT INFO

cmd="/usr/bin/fractal-explorer"

case "$1" in
  start)
    if [ -z "$(/usr/bin/pidof $cmd)" ]; then
      . /etc/profile
      $cmd &
    fi
    ;;

  stop)
    pid="$(/usr/bin/pidof $cmd)"
    if [ -n "$pid" ]; then
      kill "$pid"
    fi
    ;;

  restart)
    $0 stop
    sleep 1
    $0 start
    ;;

  *)
    echo "usage: $0 { start | stop | restart }"
    ;;
esac

exit 0
